# scripts/systemd/r3mes-miner.service
[Unit]
Description=R3MES Miner Engine
After=network.target remesd.service
Wants=network-online.target

[Service]
Type=simple
User=miner
Group=miner
WorkingDirectory=/home/miner
Environment="PATH=/home/miner/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/miner"
Environment="R3MES_NODE_GRPC_ADDR=localhost:9090"
Environment="R3MES_IPFS_ADDR=/ip4/127.0.0.1/tcp/5001"
ExecStart=/home/miner/.local/bin/r3mes-miner start --continuous
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=r3mes-miner

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (adjust based on GPU)
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target

