# scripts/systemd/r3mes-log-cleanup.timer
[Unit]
Description=R3MES Log Cleanup Timer
Requires=r3mes-log-cleanup.service

[Timer]
# Daily cleanup at 3 AM
OnCalendar=daily
OnCalendar=*-*-* 03:00:00
# Run immediately if missed
Persistent=true
# Random delay to avoid thundering herd (0-30 minutes)
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target

