# scripts/systemd/r3mes-backend.service
[Unit]
Description=R3MES Backend Inference Service
After=network.target remesd.service postgresql.service
Requires=network.target

[Service]
Type=simple
User=remes
Group=remes
WorkingDirectory=/home/remes/R3MES
Environment="PYTHONPATH=/home/remes/R3MES/backend"
Environment="R3MES_ENV=production"
Environment="BASE_MODEL_PATH=checkpoints/base_model"
Environment="CHAIN_JSON_PATH=chain.json"
# Database Configuration (PostgreSQL for production)
Environment="DATABASE_TYPE=postgresql"
Environment="DATABASE_URL=postgresql://user:password@localhost:5432/r3mes"
# Secret Management Configuration
Environment="SECRET_MANAGER_TYPE=aws"
Environment="SECRET_MANAGER_CONFIG={\"region\":\"us-east-1\"}"
# Logging Configuration
Environment="LOG_LEVEL=INFO"
Environment="ENABLE_FILE_LOGGING=true"
# Connection Pool Settings
Environment="DATABASE_POOL_MIN_SIZE=10"
Environment="DATABASE_POOL_MAX_SIZE=50"
ExecStart=/bin/bash -c 'source backend/venv/bin/activate && python run_backend.py'
Restart=always
RestartSec=5
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=r3mes-backend

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/remes/R3MES/backend/logs /home/remes/R3MES/checkpoints

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target

