# scripts/systemd/remesd.service
[Unit]
Description=R3MES Blockchain Node
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=remes
Group=remes
WorkingDirectory=/home/remes/.remesd
ExecStart=/usr/local/bin/remesd start --home /home/remes/.remesd
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=remesd

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/remes/.remesd

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment="HOME=/home/remes"
Environment="R3MES_GRPC_ADDR=0.0.0.0:9090"
Environment="R3MES_TENDERMINT_RPC_ADDR=http://0.0.0.0:26657"

[Install]
WantedBy=multi-user.target

